# wazuh_rules/100_system_rest_rules.xml
<!--
  ╔══════════════════════════════════════════════════════════════════╗
  ║  BlueTeam Challenge — System-Level & REST/JSON Rules            ║
  ║  IDs: 100700 – 100749                                           ║
  ╚══════════════════════════════════════════════════════════════════╝
  Covers:
    - REST API abuse (rate limiting, mass assignment, HTTP verb tampering)
    - JSON injection / prototype pollution
    - Linux privilege escalation
    - Reverse shell indicators
    - Crontab persistence
    - SSH brute force / key injection
    - Suspicious process execution
-->
<group name="system,rest_api,blueteam,">

  <!-- ════════════════════════════════════════════════════════════
       REST / JSON API Abuse
       ════════════════════════════════════════════════════════════ -->

  <rule id="100700" level="3">
    <match>GET |POST |PUT |DELETE |PATCH |OPTIONS </match>
    <description>Web: HTTP REST request</description>
    <group>rest_api,</group>
  </rule>

  <!-- Mass assignment / parameter pollution -->
  <rule id="100701" level="11">
    <if_sid>100700</if_sid>
    <regex>(?i)("role"\s*:\s*"admin"|"is_admin"\s*:\s*true|"isAdmin"\s*:\s*true|"admin"\s*:\s*true|"privilege"\s*:\s*|"__proto__"\s*:|"\$__proto__"\s*:|"constructor"\s*:\s*\{)</regex>
    <description>REST/JSON: Mass assignment or prototype pollution attempt</description>
    <mitre>
      <id>T1190</id>
    </mitre>
    <group>rest_api,mass_assignment,prototype_pollution,attack,</group>
  </rule>

  <!-- HTTP verb tampering -->
  <rule id="100702" level="10">
    <if_sid>100700</if_sid>
    <regex>(?i)(X-HTTP-Method-Override|X-Method-Override|X-HTTP-Method)\s*:\s*(DELETE|PUT|PATCH)</regex>
    <description>REST: HTTP verb tampering via override header</description>
    <mitre>
      <id>T1190</id>
    </mitre>
    <group>rest_api,verb_tampering,</group>
  </rule>

  <!-- IDOR / Broken Object Level Auth -->
  <rule id="100703" level="10" frequency="20" timeframe="30">
    <if_matched_sid>100700</if_matched_sid>
    <same_source_ip/>
    <match> 200 </match>
    <description>REST: High request rate with 200 OK — possible IDOR/enumeration</description>
    <mitre>
      <id>T1595.003</id>
    </mitre>
    <group>rest_api,enumeration,</group>
  </rule>

  <!-- ════════════════════════════════════════════════════════════
       System-Level: Privilege Escalation & Persistence
       ════════════════════════════════════════════════════════════ -->

  <!-- SUID/GUID abuse -->
  <rule id="100710" level="12">
    <regex>(?i)(chmod\s+[ugo+]*[sS]|chmod\s+4[0-9]{3}|chmod\s+6[0-9]{3}|find.*-perm.*-4000|find.*-perm.*-2000)</regex>
    <description>System: SUID/SGID bit set — privilege escalation risk</description>
    <mitre>
      <id>T1548.001</id>
    </mitre>
    <group>system,privilege_escalation,</group>
  </rule>

  <!-- sudo abuse -->
  <rule id="100711" level="11">
    <match>sudo: </match>
    <regex>(?i)(command not allowed|NOT in sudoers|3 incorrect password attempts|sudo su|sudo -i|sudo bash|sudo sh)</regex>
    <description>System: Suspicious sudo usage or sudo failure</description>
    <mitre>
      <id>T1548.003</id>
    </mitre>
    <group>system,privilege_escalation,</group>
  </rule>

  <!-- Reverse shell indicators -->
  <rule id="100720" level="15">
    <regex>(?i)(bash\s+-i\s+>&amp;|/dev/tcp/|nc\s+-e\s+/bin|ncat.*-e|python.*socket.*connect|perl.*socket.*connect|ruby.*socket.*connect|php.*fsockopen|mkfifo.*nc|rm\s+/tmp.*mkfifo|socat.*exec)</regex>
    <description>System: Reverse shell command detected</description>
    <mitre>
      <id>T1059</id>
    </mitre>
    <group>system,reverse_shell,attack,</group>
  </rule>

  <!-- Crontab modification (persistence) -->
  <rule id="100721" level="11">
    <regex>(?i)(crontab\s+-e|echo.*>\s*/etc/cron|tee\s+/etc/cron|>\s*/var/spool/cron|cron.*REPLACE|NEW_LINE.*cron)</regex>
    <description>System: Crontab modification detected — possible persistence</description>
    <mitre>
      <id>T1053.003</id>
    </mitre>
    <group>system,persistence,</group>
  </rule>

  <!-- SSH authorized_keys modification -->
  <rule id="100722" level="13">
    <regex>(?i)(authorized_keys|\.ssh/authorized|ssh-rsa\s+AAAA|ssh-ed25519\s+AAAA|echo.*authorized_keys|tee.*\.ssh)</regex>
    <description>System: SSH authorized_keys modified — possible backdoor</description>
    <mitre>
      <id>T1098.004</id>
    </mitre>
    <group>system,persistence,ssh,</group>
  </rule>

  <!-- New user creation -->
  <rule id="100723" level="12">
    <regex>(?i)(useradd|adduser|usermod.*-G\s+sudo|usermod.*-aG\s+sudo|passwd\s+root|chpasswd)</regex>
    <description>System: User account creation or privilege modification</description>
    <mitre>
      <id>T1136.001</id>
    </mitre>
    <group>system,persistence,account_creation,</group>
  </rule>

  <!-- Downloads from the internet to /tmp -->
  <rule id="100724" level="11">
    <regex>(?i)(wget\s+http|curl\s+-[sS]*[oO]\s+/tmp|curl\s+http.*\s+-o\s+/tmp|python.*urllib.*download|pip\s+install.*--index-url\s+http://)</regex>
    <description>System: File downloaded to /tmp — possible payload staging</description>
    <mitre>
      <id>T1105</id>
    </mitre>
    <group>system,staging,</group>
  </rule>

  <!-- /tmp execution -->
  <rule id="100725" level="12">
    <regex>(?i)(exec.*/tmp/|bash\s+/tmp/|sh\s+/tmp/|python.*\/tmp\/|chmod.*\/tmp\/.*\+x)</regex>
    <description>System: Executable run from /tmp — common malware staging area</description>
    <mitre>
      <id>T1059</id>
    </mitre>
    <group>system,malware,execution,</group>
  </rule>

  <!-- Port scanning / network recon from host -->
  <rule id="100726" level="11">
    <regex>(?i)(nmap\s|masscan\s|zmap\s|portscan|scanning\s+[0-9]+\s+ports)</regex>
    <description>System: Network port scanning tool executed</description>
    <mitre>
      <id>T1046</id>
    </mitre>
    <group>system,recon,</group>
  </rule>

  <!-- LD_PRELOAD injection -->
  <rule id="100730" level="14">
    <regex>(?i)(LD_PRELOAD\s*=|LD_LIBRARY_PATH\s*=.*\/tmp|export\s+LD_PRELOAD)</regex>
    <description>System: LD_PRELOAD injection detected — shared library hijack</description>
    <mitre>
      <id>T1574.006</id>
    </mitre>
    <group>system,defense_evasion,</group>
  </rule>

  <!-- Python pip install suspicious packages -->
  <rule id="100731" level="11">
    <regex>(?i)(pip\s+install\s+--extra-index|pip\s+install.*http://|pip3\s+install.*--index-url\s+http://)</regex>
    <description>Python: Package installed from untrusted/HTTP source — possible supply chain attack</description>
    <mitre>
      <id>T1195.001</id>
    </mitre>
    <group>system,supply_chain,</group>
  </rule>

</group>
